# automations/apple-notes-sync/com.user.apple-notes-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.apple-notes-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__REPO__/automations/apple-notes-sync/run.sh</string>
  </array>
  <!-- Once daily at 08:00 local. run.sh exits instantly when the queue is empty
       and snapshots are fresh, so most days it opens nothing; it only activates
       Notes when there's a queued card or a stale snapshot. Mac asleep at 08:00 →
       launchd runs it once on next wake (StartCalendarInterval catch-up). No
       RunAtLoad on purpose: that would fire a Notes-activating run at every
       login/reload, which is exactly what we're avoiding. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>__REPO__/automations/apple-notes-sync/.work/launchd.log</string>
  <key>StandardErrorPath</key>
  <string>__REPO__/automations/apple-notes-sync/.work/launchd.log</string>
</dict>
</plist>
