# automations/calendar-sync/com.user.calendar-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.calendar-sync</string>
  <!-- Why osascript + Terminal: launchd has no Full Disk Access, so it can't
       exec scripts under ~/Documents or reach the Chrome extension/connectors.
       Telling Terminal to run start.sh puts the session in the user TCC context. -->
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/osascript</string>
    <string>-e</string>
    <string>tell application "Terminal" to do script "__REPO__/automations/calendar-sync/start.sh"</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <!-- start.sh is idempotent and returns immediately; it is not a long-running
       process to keep alive (the tmux session is). -->
  <key>KeepAlive</key>
  <false/>
  <key>StandardOutPath</key>
  <string>__REPO__/automations/calendar-sync/.work/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>__REPO__/automations/calendar-sync/.work/launchd.err.log</string>
</dict>
</plist>
