# automations/clockify-panel/com.user.clockify-panel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.clockify-panel</string>

  <!--
    Runs the INSTALLED copy under ~/Library/Application Support (placed there by
    setup.sh), NOT the repo copy in ~/Documents. ~/Documents is TCC-protected,
    so a launchd agent's /bin/bash cannot exec scripts there without a manual
    Full Disk Access grant; ~/Library is not protected. Using `bash -lc` lets
    $HOME expand at runtime (plist paths can't expand ~ or $HOME directly).
  -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>exec "$HOME/Library/Application Support/clockify-panel/launch.sh"</string>
  </array>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/tmp/clockify-panel.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/clockify-panel.log</string>
</dict>
</plist>
