# automations/inbox-sweep-loop/com.user.inbox-sweep-loop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.inbox-sweep-loop</string>

  <!--
    Why osascript + Terminal: launchd itself does not have FDA, so it cannot
    read or exec files under ~/Documents/. Telling Terminal.app via Apple
    Event to "do script" sidesteps that: Terminal runs in the user-session
    TCC context and CAN access ~/Documents/. First time this fires, macOS may
    prompt for Automation permission ("osascript wants to control Terminal");
    approve it via System Settings → Privacy & Security → Automation.
  -->
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/osascript</string>
    <string>-e</string>
    <string>tell application "Terminal" to do script "/Users/olekorlov/Documents/GitHub/AO-Personal-OS/automations/inbox-sweep-loop/start.sh"</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/olekorlov/Documents/GitHub/AO-Personal-OS/automations/inbox-sweep-loop/.work/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/olekorlov/Documents/GitHub/AO-Personal-OS/automations/inbox-sweep-loop/.work/launchd.err.log</string>
</dict>
</plist>
